# mastra-app/mastra-app.service
[Unit]
Description=Mastra App - OCI GenAI Chat
After=network.target

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/mastra-app
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

# Restart on failure
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
MemoryLimit=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
